# deploy/demo/x402-demo@.service
[Unit]
Description=x402 demo resource workload (%i)
Documentation=https://github.com/fastnear/x402-near-facilitator
Wants=network-online.target
After=network-online.target
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
User=x402-demo-%i
Group=x402-demo-%i
WorkingDirectory=/opt/x402-demo/app
ExecStart=/usr/bin/node /opt/x402-demo/app/index.mjs
EnvironmentFile=/etc/x402-demo/%i.env

Restart=on-failure
RestartSec=5
TimeoutStartSec=30
TimeoutStopSec=15
KillSignal=SIGTERM
UMask=0077
LimitCORE=0

NoNewPrivileges=yes
PrivateDevices=yes
PrivateTmp=yes
ProtectClock=yes
ProtectControlGroups=yes
ProtectHome=yes
ProtectHostname=yes
ProtectKernelLogs=yes
ProtectKernelModules=yes
ProtectKernelTunables=yes
ProtectProc=invisible
ProtectSystem=strict
ProcSubset=pid
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictNamespaces=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
LockPersonality=yes
CapabilityBoundingSet=
AmbientCapabilities=
SystemCallArchitectures=native
SystemCallFilter=@system-service

StandardOutput=journal
StandardError=journal
SyslogIdentifier=x402-demo-%i

[Install]
WantedBy=multi-user.target
