# deploy/monitoring/x402-near-metrics.service
[Unit]
Description=x402-near operational metrics push
After=network-online.target
Wants=network-online.target
OnFailure=x402-near-alert@%n.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/x402-near-metrics.sh
# The AWS CLI needs a writable HOME for its cache; keep it in the
# unit-private /tmp so the sandbox stays read-only elsewhere.
Environment=HOME=/tmp
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
